# systemd/wsl-clean-memory.service
[Unit]
Description=Clean WSL Memory if needed

[Service]
ExecStart=/usr/local/sbin/wsl-clean-memory
