# systemd/user/wsl-forward-gpg.socket
[Unit]
Description=Forward gpg socket to Windows
PartOf=wsl-forward-gpg-all.target

[Socket]
# %t is XDG_RUNTIME_DIR
ListenStream=%t/gnupg/S.gpg-agent
SocketMode=0600
DirectoryMode=0700
Accept=yes

[Install]
WantedBy=sockets.target
