# systemd/user/wsl-test.service
[Unit]
Description=Test service

[Service]
ExecStart=%h/.local/lib/systemd/wsl-test

[Install]
WantedBy=default.target
