# systemd/user/wsl-forward-gpg.service
[Unit]
Description=Forward gpg to Windows

[Service]
ExecStart=%h/.local/lib/systemd/wsl-forward-gpg

[Install]
WantedBy=default.target
